# prop-voter.service
[Unit]
Description=Prop-Voter - Cosmos Governance Bot
Documentation=https://github.com/chalabi2/prop-voter
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/path/to/your/prop-voter
ExecStart=/path/to/your/prop-voter/prop-voter -config config.yaml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10

# Basic security settings
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/go/bin
Environment=LOG_LEVEL=info

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prop-voter

[Install]
WantedBy=multi-user.target